# shared/etc/profile.d/send_notification.sh
/usr/local/sbin/shell-login.sh | mail -s "[ssh] Connexion de ${USER} sur ${HOSTNAME}" -r mailer@biodiversite-aura.fr adminsys@biodiversite-aura.fr
